# 7.1/config/supervisor.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisord/supervisord.log

[program:php-fpm]
command=php-fpm
